# tests/networks/docker-localnet/supervisord.conf
[supervisord]
logfile = /tmp/supervisord.log
logfile_maxbytes = 50MB
logfile_backups=10
logLevel = error
pidfile = /tmp/supervisord.pid
nodaemon = true
minfds = 1024
minprocs = 200
umask = 022
user = cheqd
identifier = supervisor
directory = /tmp
nocleanup = true
childlogdir = /tmp
strip_ansi = false


[program:validator-0]
command=cheqd-noded start --home "/home/cheqd/network-config/validator-0"
directory=/home/cheqd
stdout_logfile=/tmp/validator-0.log
stderr_logfile=/tmp/validator-0.log

[program:validator-1]
command=cheqd-noded start --home "/home/cheqd/network-config/validator-1"
directory=/home/cheqd
stdout_logfile=/tmp/validator-1.log
stderr_logfile=/tmp/validator-1.log

[program:validator-0-log]
command=tail -f /tmp/validator-0.log
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:validator-1-log]
command=tail -f /tmp/validator-1.log
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
